# curator4-990829104.conf
[program:curator4-990829104]
command=python3 /home/curator4/curator4-990829104/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-990829104/test.err.log
stdout_logfile=/home/curator4/curator4-990829104/test.out.log